# target/linux/mtmips/base-files/etc/hotplug.d/iface/00-switch
#!/bin/sh

[ "$ACTION" = ifup -a "$INTERFACE" = lan ] || exit 0

. /lib/functions.sh
. /lib/network/switch.sh

logger -t switch "reset switch($INTERFACE)"
reset_switch_vlan
