# target/linux/mtmips/base-files/etc/hotplug.d/net/01-rps
#!/bin/sh

if [ "$ACTION" = "add" ]; then
	case "$DEVICENAME" in
	"eth2" | "eth3" | "ra0" | "rai0")
		smp.sh
		;;
	esac
fi
